# monitor-agent/systemd/rybbit-monitor-agent.service
[Unit]
Description=GHOST ID Monitor Agent
Documentation=https://github.com/ghost-id/monitor-agent
After=network.target

[Service]
Type=simple
User=ghost-id
Group=ghost-id
WorkingDirectory=/opt/ghost-id-monitor-agent
ExecStart=/usr/bin/node /opt/ghost-id-monitor-agent/dist/index.js
Restart=always
RestartSec=10

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/opt/ghost-id-monitor-agent/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ghost-id-monitor-agent/logs

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghost-id-monitor-agent

[Install]
WantedBy=multi-user.target